# chatroom-websocket.conf
[program:chatroom-websocket]
command=/usr/bin/python3.6 chat.py
directory=/var/www/chatroom-websocket
autostart=true
autorestart=true
